# python/rq/rq_task.ini
[program:rqtask]
directory = /Application/apid-1.5
command = /Application/apid/bin/rqworker task
process_name = rqtask-%(process_num)s
numprocs = 3
autostart = True
user=root
stopsignal=TERM
autorestart=true
stopasgroup=true
startsecs=10
startretries=999
stdout_logfile=/var/log/apid/rqtask.out
stderr_logfile=/var/log/apid/rqtask.err